# firmware/openclaw.service
[Unit]
Description=OpenClaw Daemon (Conda Isolated)
# 1. 强制要求网络在线
Wants=network-online.target
# 2. 在网络完全在线之后才启动
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/home/pi/Moving-Claw/start_openclaw.sh
# 3. 增加每次重启前的冷却时间，防止疯狂报错刷屏
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target